# src/config/upstart/icebox.conf
description "IceBox service"
author "ZeroC.com"

start on runlevel [234]
stop on runlevel [0156]

chdir /home/andreynech/projects/a20/src/config
exec icebox --Ice.Config=icebox.config

respawn
